# docker/crontab
# V2EX Max Helper Cron Tasks (Docker)
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# 每天 UTC 01:10 执行签到
10 1 * * * root . /etc/environment && cd /app/checkin && /usr/bin/node v2ex-checkin.js >> /data/v2ex-checkin.log 2>&1

# 每天 UTC 0、6、12、18 保活 (Ping)
0 0,6,12,18 * * * root . /etc/environment && cd /app/checkin && /usr/bin/node v2ex-checkin.js --ping >> /data/v2ex-checkin.log 2>&1

# 每天 UTC 01:15 执行阅读
# 使用 xvfb-run 提供虚拟桌面，解决 Headless 检测问题
15 1 * * * root . /etc/environment && cd /app/reader && xvfb-run -a /usr/bin/node main.js >> /data/v2ex-reader.log 2>&1

# 每天 UTC 05:00 清理一次遗留的 xvfb 锁和 chromium 进程（防止僵尸进程积累）
0 5 * * * root pkill -f chromium; rm -f /tmp/.X99-lock /tmp/.X11-unix/X99 >> /dev/null 2>&1

# 下面必须留一空行

